# systemd/pr_sniper.service
[Unit]
Description=GitHub PR Watcher Bot
After=network.target

[Service]
Type=simple
ExecStart=/opt/pr_sniper/venv/bin/python /opt/pr_sniper/pr_sniper.py
WorkingDirectory=/opt/pr_sniper
Restart=always
RestartSec=5
User=YOUR_USERNAME
Group=YOUR_USERNAME
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

